# swarm-service.conf
description "Manage Browsh Swarm nodes"

start on filesystem or runlevel [2345]
stop on shutdown

exec /usr/local/bin/swarm-manager.sh
